# enhanced-betting-ai.service
[Unit]
Description=Enhanced Sports Betting AI Monitoring
After=network.target

[Service]
Type=simple
User=reid
WorkingDirectory=/home/reid/Desktop/parleyapp
ExecStart=/home/reid/.pyenv/versions/3.12.11/bin/python /home/reid/Desktop/parleyapp/start_monitoring.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
